# overlay/etc/periodic/15min/wp-cron
#!/usr/bin/with-contenv bash
[[ "${USE_SYSTEM_CRON}" = "1" ]] && wp-cron
