# overlay/etc/periodic/15min/wp-cron
#!/bin/bash

wget -O - http://localhost/wp-cron.php
