# overlay/etc/periodic/15min/wp-cron
#!/bin/bash

echo " - executing WordPress cron..."
wget -O - http://localhost/wp-cron.php
echo " - done."
