# overlay/etc/periodic/15min/wp-cron
#!/usr/bin/with-contenv bash

set -euo pipefail


#======================================================================================================================
# Call wp-cron executable
#======================================================================================================================

[[ "${USE_SYSTEM_CRON:-0}" = "1" ]] && wp-cron
