# overlay/etc/periodic/15min/wp-cron
#!/usr/bin/with-contenv bash

_echo "executing WordPress cron..."
wget -O - http://localhost/wp-cron.php
_done
